# template/ic-r232-1.template.rules

ACTION=="add|remove", SUBSYSTEM=="tty", DEVPATH=="*/tty/{{tty_b}}", RUN+="/opt/unipi/tools/uhelper rs232.{{slot}} tty {{slot}} $name"
